# deploy/systemd/opencode-server.service.d/10-disable-project-config.conf
[Service]
Environment=OPENCODE_DISABLE_PROJECT_CONFIG=1
